# btrfs-snapshot@.timer
[Unit]
Description=Hourly snapshot of %f

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
