# btrfs-snapshot-cleanup@.timer
[Unit]
Description=Daily snapshot cleanup of %f

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
