# btrfs-snapshot@.service.template
[Unit]
Description=Take a snapshot of %f

[Service]
Slice=background.slice
Type=oneshot
ExecStart=PREFIX/bin/btrfs-snapshot %f

[Install]
Also=btrfs-snapshot@%i.timer
Also=btrfs-snapshot-cleanup@%i.service
