# btrfs-snapshot@.timer
[Unit]
Description=Hourly snapshot of %I

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
